# config/system_trader.service
[Unit]
Description=System Trader Services
After=network.target postgresql.service redis-server.service

[Service]
Type=forking
User=ubuntu
WorkingDirectory=/home/ubuntu/system_trader
ExecStart=/home/ubuntu/system_trader/scripts/start_services.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target